# templates/supervisor/closure.conf
{% set programs = [] %}

{% for transform in transforms %}
{% do programs.append(transform.name) %}
[program:{{transform.name}}]
command=/Users/arpitbh/workspace/stacs/ui/venv/bin/python /Users/arpitbh/workspace/stacs/ui/executors/transform-{{transform.name}}-{{transform.input_source.type}}-{{transform.input_source.meta.format}}-{{transform.output_target.type}}-{{transform.output_target.meta.format}}.py {{transform._id}}
autostart=False
{% endfor %}

{% for puller in pullers %}
{% do programs.append(puller.name) %}
[program:{{puller.name}}]
command=/Users/arpitbh/workspace/stacs/ui/venv/bin/python /Users/arpitbh/workspace/stacs/ui/executors/puller-{{puller.input_source.type}}-{{puller.input_source.meta.format}}-{{puller.pipeline.type}}-{{puller.pipeline.meta.format}}.py {{puller._id}}
autostart=False
{% endfor %}

{% for pusher in pushers %}
{% do programs.append(pusher.name) %}
[program:{{pusher.name}}]
command=/Users/arpitbh/workspace/stacs/ui/venv/bin/python /Users/arpitbh/workspace/stacs/ui/executors/pusher-{{pusher.pipeline.type}}-{{pusher.pipeline.meta.format}}-{{pusher.output_target.type}}-{{pusher.output_target.meta.format}}.py {{pusher._id}}
autostart=False
{% endfor %}

[group:{{closure.name}}]
programs={{programs|join(',')}}
